# roles/backup/templates/vault-backup.service.j2
[Unit]
Description=Vault Backup Service
After=vault.service
Wants=vault.service

[Service]
Type=oneshot
User={{ backup_user }}
Group={{ backup_group }}
Environment=VAULT_ADDR=https://127.0.0.1:8200
Environment=VAULT_SKIP_VERIFY=true
Environment=RESTIC_PASSWORD_FILE={{ backup_config_dir }}/restic-password
Environment=RESTIC_REPOSITORY={{ restic_repository_nas }}
ExecStart=/usr/local/bin/vault-backup.sh
StandardOutput=journal
StandardError=journal
SyslogIdentifier=vault-backup

# Resource limits
CPUQuota={{ backup_cpu_limit }}
MemoryLimit={{ backup_memory_limit }}

# Security
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths={{ backup_root_dir }}
ReadWritePaths={{ backup_config_dir }}
ReadWritePaths={{ backup_log_dir }}

[Install]
WantedBy=multi-user.target